# root/etc/keybr-backup-database.timer
[Unit]
Description=keybr.com database backup service timer

[Timer]
OnCalendar=*-*-* 00:30:00

[Install]
WantedBy=multi-user.target
