# root/etc/keybr-remove-sessions.timer
[Unit]
Description=keybr.com stale session removal service timer

[Timer]
OnCalendar=*-*-* *:15:00

[Install]
WantedBy=multi-user.target
